# util/pigpiod.service
[Unit]
Description=Pigpio daemon

[Service]
Type=simple
ExecStart=/usr/bin/pigpiod -g

[Install]
WantedBy=multi-user.target
